# home_encryption/service/fallback-to-cli.service
[Unit]
Description=Fall back to multi-user target (CLI)

[Service]
Type=oneshot
ExecStart=/bin/systemctl isolate multi-user.target

